# docker/supervisor.conf
[supervisord]
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:hugo]
command=/usr/local/bin/hugo server --bind 0.0.0.0 --buildDrafts
directory=/mnt/data

[program:webpack]
command=/mnt/data/themes/blackburn/node_modules/.bin/webpack --watch
directory=/mnt/data/themes/blackburn

[program:gulp]
command=/mnt/data/node_modules/.bin/gulp thumbnails mini_png mini_jpg watch
directory=/mnt/data

[unix_http_server]
file = /tmp/supervisor.sock

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface